# systemd/cups-keepalive.service
[Unit]
Description=CUPS Keepalive Print Worker
After=network.target cups.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cups-keepalive-worker run
